# contrib/xdg-desktop-portal.shiny.service.in
[Unit]
Description=Portal service for shiny-shell
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.shiny
ExecStart=@libexecdir@/xdg-desktop-portal-shiny
Restart=on-failure
Slice=session.slice
